# kubernetes/samples/kubelet.service.tmpl
[Service]
ExecStart=
ExecStart=/usr/bin/kubelet --kubeconfig=/etc/kubernetes/kubelet.conf --config=/var/lib/kubelet/config.yaml --node-ip='${IPV4},${IPV6}'
